# deploy/ugetfirst-health.service
[Unit]
Description=UGetFirst engine HTTP health check (/health on port 8080)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ugetfirst
Group=ugetfirst
WorkingDirectory=/home/ugetfirst/UGetFirst_engine
EnvironmentFile=/home/ugetfirst/UGetFirst_engine/.env
ExecStart=/home/ugetfirst/UGetFirst_engine/.venv/bin/python health_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
