# deploy/ugetfirst-sync-prod-to-dev.timer
[Unit]
Description=Weekly PROD → DEV Supabase sync

[Timer]
# Sunday at midnight UTC
OnCalendar=Sun *-*-* 00:00:00 UTC
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
